# debian/config/init.d/kanata
#!/sbin/openrc-run

command="/home/f1/.cargo/bin/kanata"
#command_args="--config=/home/<username>/.config/kanata/kanata.kbd"

command_background=true
pidfile="/run/${RC_SVCNAME}.pid"

command_user="f1"
