# systemd/user/dbus_monitor_date_session.service
[Unit]
Description=i3_daemon
PartOf=graphical-session.target

[Service]
ExecStart=%h/dot/bin/dbus_monitor_date --session
Restart=always
