# shared/config/systemd/user/i3_daemon.service
[Unit]
Description=i3_daemon
PartOf=graphical-session.target

[Service]
ExecStart=i3_daemon
Restart=always

[Install]
WantedBy=xsession.target
